# sgs.service
[Unit]
Description=Smart Group Server
After=systemd-user-session.service

[Service]
Type=simple
ExecStart=/usr/local/bin/sgs /usr/local/etc/sgs.cfg
